# to-deploy/autostart/rvgl.desktop
[Desktop Entry]
Type=Application
Name=rvgl
Exec=/home/rvgl/rvgl.sh
